# dock/depage/images/php-5.6/etc/service/php-fpm/run
#!/bin/sh

exec 2>&1
exec /usr/bin/php-fpm5.6 -y /etc/php/fpm-php5.6/php-fpm.conf -F
